# dot_config/crontab
*/15 * * * * ~/.local/bin/update-neovim.sh
@daily chezmoi update
